# customize/files/logrotate.cron
00 00 * * * /usr/sbin/logrotate /etc/logrotate.conf
